# vgresources/teampostgresql.conf
description "PostgreSQL web gui"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /srv/teampostgresql
setuid postgres

exec /srv/teampostgresql/teampostgresql-run.sh
